# crontab
#crontab
* * * * *  /venv/bin/python3  /app/app.py
#end of line