# crontab
#crontab
* * * * *  /venv/bin/python3  /CronJob/app.py
#end of line